# support/systemd/vepc-pgwd.service.in
[Unit]
Description=VEPC P-GW Daemon
After=networking.service

[Service]
Type=simple

User=vepc
Group=vepc

RuntimeDirectory=vepc-pgwd
PIDFile=@LIB_DIR@/run/vepc-pgwd/pid

Restart=always
ExecStart=@BIN_DIR@/vepc-pgwd -f @SYSCONF_DIR@/vepc/pgw.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
