# support/systemd/vepc-mmed.service.in
[Unit]
Description=VEPC MME Daemon
After=networking.service

[Service]
Type=simple

User=vepc
Group=vepc

RuntimeDirectory=vepc-mmed
PIDFile=@LIB_DIR@/run/vepc-mmed/pid

Restart=always
ExecStart=@BIN_DIR@/vepc-mmed -f @SYSCONF_DIR@/vepc/mme.conf
RestartSec=2
RestartPreventExitStatus=1

[Install]
WantedBy=multi-user.target
