# support/root/etc/cron.weekly/vepc_log_clean.sh
#!/bin/bash
echo "++++++++" > /root/app/vepc/log/vepc.log
